# code/src/service/obsolete/collectorwd.service
[Unit]
Description=Collector Web Service
After=mysqld.service
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=gvalera
ExecStart=/usr/local/bin/collectorwd

[Install]
WantedBy=multi-user.target
